# installer/cdrom/assets/etc/systemd/system/plymouth-quit.service.d/override.conf
[Unit]
After=kiosk-installer.service
Requires=kiosk-installer.service
