# overlay/etc/local.d/10-setup-alpine.start
#! /bin/sh

set -o errexit
set -o nounset

# Close standard input.
exec 0<&-

setup-alpine -ef /etc/auto-setup-alpine/answers
